# k8s-emulator/utils/tc_mod.service.sample
[Unit]
Description=Activate automated tc
[Service]
ExecStart=/vagrant/utils/tc_mod.sh